# etc/backups/crontab
# Execute the generator script every day at 2:30 AM
30 2 * * * /bin/sh /bin/generator
# Empty line